# bin/data/scripts/com.ofappRelaunch.openFrameworks.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>KeepAlive</key>
	<dict>
		<key>SuccessfulExit</key>
		<false/>
	</dict>
	<key>Label</key>
	<string>com.ofappRelaunch.openFrameworks</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/masterstudio/Documents/00-App/OF/9_8/apps/myApps/led_grid_v2/bin/led_grid_v2.app</string>
	</array>
</dict>
</plist>
